# services/backend/schedule.cron
* * * * * python $BACKEND_PATH/manage.py runcrons >> /var/log/cron.log 2>&1
